# .config/anacrontab/on_ac
# /usr/sbin/anacron -S "$XDG_DATA_HOME/anacron" -t "$XDG_CONFIG_HOME/anacrontab" -s
SHELL=/bin/bash
USER=KONFEKT
HOME=/home/KONFEKT
# See http://unix.stackexchange.com/questions/138577/how-to-guarantee-availability-of-bash-env/138609#138609
BASH_ENV=/home/KONFEKT/.bash_profile
# or:
# @daily       0      daily    bash -l -c 'run-parts.sh $HOME/.config/anacron.daily'
XDG_CONFIG_HOME=/home/KONFEKT/.config
XDG_CACHE_HOME=/home/KONFEKT/.cache
XDG_DATA_HOME=/home/KONFEKT/.local/share
PATH=/usr/local/bin:/usr/bin:/usr/sbin:/bin:/sbin:/home/KONFEKT/bin:/home/KONFEKT/.local/bin
# DISPLAY=:0
MAILTO=KONFEKT
RANDOM_DELAY=5
# START_HOURS_RANGE=10-16

# period  delay  job-identifier  command
@daily    0   daily.on_ac  cronic chrt --idle 0 ionice -c2 -n7 run-parts.sh -v $XDG_CONFIG_HOME/anacron.daily/on_ac
@weekly   5   weekly.on_ac       cronic chrt --idle 0 ionice -c2 -n7 run-parts.sh -v $XDG_CONFIG_HOME/anacron.weekly/on_ac
@monthly  30  monthly.on_ac      cronic chrt --idle 0 ionice -c2 -n7 run-parts.sh -v $XDG_CONFIG_HOME/anacron.monthly/on_ac

# ex:ft=crontab 
